# runtime-images/ubuntu-24.04/gnome/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:runtime-setup]
command=/usr/local/bin/prepare-runtime.sh
priority=1
autorestart=false
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dbus-system]
command=/usr/bin/dbus-daemon --system --nofork
priority=5
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:systemd-logind]
command=/lib/systemd/systemd-logind
priority=6
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vncserver]
command=/usr/local/bin/start-vnc.sh
priority=10
autorestart=true
user=desktop
environment=HOME="/home/desktop",USER="desktop",DISPLAY=":1",XDG_RUNTIME_DIR="/run/user/1001"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ 6080 localhost:5901
priority=20
autorestart=true
user=desktop
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
